# vendor/github.com/hashicorp/nomad/dist/upstart/nomad.conf
start on (filesystem and net-device-up IFACE=lo)
stop on runlevel [!2345]

exec /usr/bin/nomad agent -config /etc/nomad